# Discord Agent/launchd/com.moonshots.firstmate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.moonshots.firstmate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/citadel/Projects/MoonshotsMates/Discord Agent/.venv/bin/python3</string>
    <string>/Users/citadel/Projects/MoonshotsMates/Discord Agent/bot.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/citadel/Projects/MoonshotsMates/Discord Agent</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/citadel/Projects/MoonshotsMates/Discord Agent/data/bot.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/citadel/Projects/MoonshotsMates/Discord Agent/data/bot.log</string>
</dict>
</plist>
